# final/capstone_monitor.service
	# capstone_monitor.service — Real-Time Performance Monitoring Daemon
#
# INSTALL:
#   sudo make install                          # builds and copies to /opt/capstone_monitor/
#   sudo cp capstone_monitor.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now capstone_monitor
#
# VIEW OUTPUT:
#   journalctl -u capstone_monitor -f
#
# CHECK STATUS:
#   systemctl status capstone_monitor

[Unit]
Description=Capstone Real-Time Performance Monitoring Daemon
Documentation=https://github.com/group4/capstone_monitor
# Start after the kernel's BPF subsystem and filesystem are available.
# network.target not required — daemon is local-only.
After=local-fs.target
# If the classifier or ghOSt agent are separate services, uncomment:
# Before=capstone_classifier.service capstone_ghost_agent.service
# Don't restart more than 5 times in 60 seconds — prevents crash loop thrash.
StartLimitBurst=5
StartLimitIntervalSec=60	

[Service]
# Item 4: Type=notify lets systemd wait for sd_notify_ready() before marking
# the unit active.  This means any service declared After=capstone_monitor.service
# (e.g. capstone_classifier.service, capstone_ghost_agent.service) will not
# start until the ring buffers are mapped and the eBPF programs are attached —
# guaranteeing consumers never try to attach before the shm regions exist.
# The daemon calls sd_notify("READY=1") via an inline implementation that
# writes directly to NOTIFY_SOCKET, so libsystemd is not required.
Type=notify
ExecStart=/opt/capstone_monitor/monitor_final /opt/capstone_monitor/bpf/ctx_switch.bpf.o --all --top 50 --interval 50

# Working directory for any relative path resolution
WorkingDirectory=/opt/capstone_monitor

# Run as root so perf_event_open and bpf() syscalls succeed.
# Alternative: run as non-root with the capabilities below — requires
# kernel 5.9+ and sysctl kernel.perf_event_paranoid <= 2.
User=root
Group=root

# Capabilities needed:
#   CAP_PERFMON   — perf_event_open() for hardware counters (kernel 5.9+)
#                   On older kernels this was part of CAP_SYS_ADMIN.
#   CAP_BPF       — load and attach eBPF programs (kernel 5.8+)
#                   On older kernels this was part of CAP_SYS_ADMIN.
#   CAP_SYS_ADMIN — fallback for kernels < 5.8 that don't split these out.
#                   Remove if running kernel >= 5.9 to reduce attack surface.
#   CAP_DAC_READ_SEARCH — read /proc/<pid>/ entries owned by other users.
AmbientCapabilities=CAP_PERFMON CAP_BPF CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_PERFMON CAP_BPF CAP_SYS_ADMIN CAP_DAC_READ_SEARCH

# File descriptor limit.
# Calculation: MAX_PIDS * (3 fds aggregate + MAX_CORES * 3 fds per-core)
#              = 256 * (3 + 64*3) = 256 * 195 = 49920
# Round up to next power-of-2 with headroom.
LimitNOFILE=65536

# Memory lock limit for mmap'd shared memory ring buffers.
# Two RingBuffer regions: 2 * sizeof(RingBuffer)
# sizeof(RingBuffer) = RB_CAPACITY(32768) * sizeof(SmoothedMetrics)(48) = ~1.5MB
# Set to 8MB for headroom.
LimitMEMLOCK=8388608

# Restart policy: restart on crash/OOM, but not on clean exit (e.g. SIGTERM).
Restart=on-failure
RestartSec=2s


# Log stdout/stderr to the journal.
# View with: journalctl -u capstone_monitor -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=capstone_monitor

# Create the CSV log directory if it doesn't exist.
RuntimeDirectory=capstone_monitor
LogsDirectory=capstone_monitor
LogsDirectoryMode=0755

# Security hardening — restrict what the daemon can do beyond its capabilities.
NoNewPrivileges=yes
ProtectSystem=strict
# Allow writing to /dev/shm (POSIX shared memory for ring buffers)
# and /var/log/capstone_monitor (CSV output).
ReadWritePaths=/dev/shm /var/log/capstone_monitor
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
